# guest/hadoop/files/etc/rc.d/rc.local
#!/bin/bash

if [ -f /vagrant/.vagrant/.ssh/id_rsa ]; then 
	cp /vagrant/.vagrant/.ssh/* ~/.ssh/ 
else
	ssh-keygen -t rsa -P '' -f ~/.ssh/id_rsa
	cat ~/.ssh/id_rsa.pub >> ~/.ssh/authorized_keys 
	chmod 0600 ~/.ssh/authorized_keys
	[ -d /vagrant/.vagrant/.ssh ] || mkdir -p /vagrant/.vagrant/.ssh
	cp ~/.ssh/* vagrant/.vagrant/.ssh/
fi

$HADOOP_HOME/bin/hdfs namenode -format
$HADOOP_HOME/sbin/start-dfs.sh

