# guest/hadoop/files/etc/rc.d/rc.local
#!/bin/bash -lex

if [ -d /vagrant/.vagrant/.ssh ]; then 
	cp -r /vagrant/.vagrant/.ssh ~/.ssh 
else
	ssh-keygen -t rsa -P '' -f ~/.ssh/id_rsa
	cat ~/.ssh/id_rsa.pub >> ~/.ssh/authorized_keys 
	chmod 0600 ~/.ssh/authorized_keys
	cp -r ~/.ssh vagrant/.vagrant/.ssh
fi

cat <<_EOT_>> /etc/resolv.conf
nameserver $(getent hosts dpsrv-bind-1|awk '{ print $1 }')
search dpsrv.me
_EOT_

[ -d /tmp/hadoop-root/dfs/name ] || $HADOOP_HOME/bin/hdfs namenode -format
$HADOOP_HOME/sbin/start-dfs.sh

